# misc/systemd/nix-daemon.service.in
[Unit]
Description=Nix Daemon
Documentation=man:nix-daemon https://nixos.org/manual
RequiresMountsFor=@storedir@
RequiresMountsFor=@localstatedir@
RequiresMountsFor=@localstatedir@/nix/db
ConditionPathIsReadWrite=@localstatedir@/nix/daemon-socket

[Service]
ExecStart=@@bindir@/nix-daemon nix-daemon --daemon
# The daemon forks per-client workers and build processes. Stop/restart must
# terminate the whole service cgroup so stale workers cannot survive a restart
# and keep store locks alive.
KillMode=control-group
LimitNOFILE=1048576
TasksMax=1048576
Delegate=yes

[Install]
WantedBy=multi-user.target
